# crypto-trading.service
# /etc/systemd/system/crypto-trading.service
# Servicio SystemD para ejecutar sin Docker (alternativa)

[Unit]
Description=Crypto Trading Simulator Backend
After=network.target mysql.service
Requires=mysql.service

[Service]
Type=simple
User=crypto-app
Group=crypto-app

# Variables de entorno
Environment="SPRING_PROFILES_ACTIVE=production"
Environment="SERVER_PORT=8080"
EnvironmentFile=/opt/crypto-trading/app.env

# Directorio de trabajo
WorkingDirectory=/opt/crypto-trading

# Comando de ejecución
ExecStart=/usr/bin/java -Xmx512m -Xms256m \
    -Dspring.profiles.active=production \
    -Dserver.port=8080 \
    -jar /opt/crypto-trading/app.jar

# Reinicio automático
Restart=always
RestartSec=10

# Logs
StandardOutput=append:/opt/crypto-trading/logs/app.log
StandardError=append:/opt/crypto-trading/logs/error.log

# Límites
LimitNOFILE=65536
LimitNPROC=32768

# Timeout
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
